# addon/rootfs/etc/services.d/amazonparent/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# S6 service finish script for Amazon Parent Dashboard Auth
# ==============================================================================

bashio::log.info "Amazon Parent Dashboard Auth Service stopped"
