# addon/rootfs/etc/services.d/amazonparent/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# S6 service runner for Amazon Parent Dashboard Auth
# ==============================================================================

exec /usr/local/bin/run.sh
